# tex/services/update-texlive.timer
[Unit]
Description="Update TeX Live"

[Timer]
OnBootSec=15min
OnUnitActiveSec=8h
Unit=update-texlive.service

[Install]
WantedBy=timers.target
