# etc/systemd/system/bootc-fetch-apply-updates.timer.d/override.conf
[Timer]
OnBootSec=
OnUnitInactiveSec=
RandomizedDelaySec=

OnCalendar=01:15:00 America/Edmonton
